# scripts/launchd/com.pepper.agent-heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- NOTE: Run `make agents-install` to generate the installed version
     with absolute paths filled in from your repo location. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pepper.agent-heartbeat</string>
    <key>ProgramArguments</key>
    <array>
        <string>REPO_ROOT/scripts/agent-heartbeat.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>7200</integer>
    <key>StandardOutPath</key>
    <string>REPO_ROOT/build/logs/heartbeat.log</string>
    <key>StandardErrorPath</key>
    <string>REPO_ROOT/build/logs/heartbeat.log</string>
    <key>WorkingDirectory</key>
    <string>REPO_ROOT</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
